# qr3_73.conf
[program:qr3_73]
command=python3 /home/curator3/qr3_73/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_73/test.err.log
stdout_logfile=/home/curator3/qr3_73/test.err.log
